# image/rootfs/usr/lib/systemd/system/nix-daemon.service.d/10-bluecat-atomic.conf
[Unit]
Requires=nix.mount bluecat-nix-setup.service
After=nix.mount bluecat-nix-setup.service
