# system_files/usr/lib/systemd/system/bluecat-install-brave.service
[Unit]
Description=Install the Brave browser Flatpak on first boot
Documentation=https://github.com/echocat/bluecat
# Run until it succeeds once; the marker is created on success below. If the
# first boot has no network, the unit fails and is retried (see Restart=), and
# it is attempted again on subsequent boots because the marker is absent.
ConditionPathExists=!/var/lib/flatpak/.bluecat-brave-installed
# Needs the Flathub remote and a working network.
After=bluecat-add-flathub.service network-online.target
Wants=network-online.target
Requires=bluecat-add-flathub.service
# Bound the in-boot retries so we do not loop forever on a persistent failure.
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=oneshot
# Install Brave system-wide so every user (including newly created ones) has it.
ExecStart=/usr/bin/flatpak install --system --noninteractive --assumeyes flathub com.brave.Browser
ExecStartPost=/usr/bin/touch /var/lib/flatpak/.bluecat-brave-installed
# Retry within this boot if it fails (e.g. transient network), then give up and
# let the next boot try again (marker is only written on success).
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
